# scripts/crontab
# every minute, run the laravel scheduler
* * * * *  vagrant  php /home/vagrant/ats/artisan scheduled:run
